# deploy/systemd/ipcheck-collector.service
[Unit]
Description=ipcheck OpenResty log collector
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/ubuntu/personal/workSpace/ipcheck/apps/api
Environment=NODE_ENV=production
Environment=IPCHECK_DB_PATH=/var/lib/ipcheck/ipcheck.db
Environment=IPCHECK_TAIL_LINES=500000
ExecStart=/home/ubuntu/.nvm/versions/node/v24.14.1/bin/node dist/collect.js
User=ubuntu
Group=ubuntu
ReadWritePaths=/var/lib/ipcheck /var/log/ipcheck
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
CapabilityBoundingSet=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
